# systemd/topology.service
[Unit]
Description=Regenerate homelab topology map
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# adjust user + path to your checkout
User=joe
WorkingDirectory=/opt/topologygenerator
ExecStart=/usr/bin/python3 make_network_topology.py --config config.yaml
# arpscan needs raw sockets; grant just that instead of full root if desired
AmbientCapabilities=CAP_NET_RAW
